# com.serverpulse.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.serverpulse.app</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/streamlit</string>
        <string>run</string>
        <string>/Users/daijun/.openclaw/workspace/projects/ServerPulse/dashboard.py</string>
        <string>--server.port</string>
        <string>8080</string>
        <string>--server.address</string>
        <string>0.0.0.0</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/serverpulse.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/serverpulse-error.log</string>
</dict>
</plist>
